# layout/Library/LaunchDaemons/com.t7rutweaks.fakepicturebypass.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.t7rutweaks.fakepicturebypass</string>
	<key>MachServices</key>
	<dict>
		<key>com.t7rutweaks.fakepicturebypass</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/FakePictureBypass</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
